# services/binary_10min.service
[Unit]
Description=OPTION 10 MIN BINARY
After=network.target

[Service]
User=vova
Group=vova

WorkingDirectory=/home/vova/Options/Binary/10min/
Environment="PYTHONPATH=/home/vova/Options/Binary/10min/"
ExecStart=/home/vova/Options/Binary/10min/venv/bin/python3.11 /home/vova/Options/Binary/10min/bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target